# systemd/homers-lotr-backup.service
[Unit]
Description=Homer's LotR server state backup (vault + databases + configs) to OneDrive
Documentation=file:///var/home/james/GIT/nwn_homers_lotr/bin/backup-homers-lotr

[Service]
Type=oneshot
# Self-gates on a 24h sentinel, so firing more often than daily is harmless.
ExecStart=/var/home/james/GIT/nwn_homers_lotr/bin/backup-homers-lotr
Nice=10
IOSchedulingClass=idle
