# systemd/homers-lotr-backup.timer
[Unit]
Description=Daily Homer's LotR server state backup
Documentation=file:///var/home/james/GIT/nwn_homers_lotr/bin/backup-homers-lotr

[Timer]
# Fire daily. Persistent=true runs a missed backup at next login/boot if the
# machine was off at the scheduled time (the host is not guaranteed up daily).
OnCalendar=daily
Persistent=true
# Spread the run a little so it doesn't fight a fresh login.
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
