# systemd/homers-lotr-wiki-publish.service
[Unit]
Description=Full Homer's LotR wiki republish + push (runs once per boot)
Documentation=file:///var/home/james/GIT/nwn_homers_lotr/bin/refresh-homers-lotr-wiki
# After the daily reboot the box comes back up and republishes the WHOLE wiki so
# creature-index/detail kill counts (not just the activity pages) reflect the
# previous day's play. Idempotent: no commit when nothing changed.
After=default.target

[Service]
Type=oneshot
# Generous timeout: a full regen takes ~1-2 min, plus up to ~2.5 min of push
# retries while the network finishes coming up after boot.
TimeoutStartSec=1200
Environment=PATH=/var/home/james/GIT/nwn_manager/bin:/home/james/.nimble/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/var/home/james/GIT/nwn_homers_lotr/bin/refresh-homers-lotr-wiki --publish
Nice=10

[Install]
WantedBy=default.target
